# supervisord.conf
[supervisord]
nodaemon=true

[program:rate-my-ride-users-api]
command=node /rate-my-ride-users-api/server.js