# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log

[program:rate-my-ride-users-api]
directory=/rate-my-ride-users-api
command=node server.js